# deploy/telegram-monitor.service
[Unit]
Description=Telegram Comment Monitor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/monitorshik-latest
EnvironmentFile=/home/ubuntu/monitorshik-latest/.env
ExecStart=/home/ubuntu/monitorshik-latest/venv/bin/python /home/ubuntu/monitorshik-latest/worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

